# dot_config/systemd/user/chezmoi-update.service
[Unit]
Description=Update dotfiles
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=chezmoi update --force
