# dot_config/systemd/user/restic.timer
[Unit]
Description=Hourly backups

[Timer]
OnStartupSec=15m
OnUnitInactiveSec=1h

[Install]
WantedBy=timers.target
