# dot_config/systemd/user/steamdeck-backup.service
[Unit]
Description=Back up saves hourly

[Timer]
OnStartupSec=15m
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
